# roles/build_node_in_docker/templates/reposync.service.j2
[Unit]
Description=Synchronize AlmaLinux local repos through DNF reposync

[Service]
Type=simple
ExecStart=/usr/local/bin/{{ update_unit_name }}.sh

[Install]
WantedBy=multi-user.target
